# docs/systemd/ckan.service
[Unit]

Description=CKAN service

After=docker.service

[Service]

ExecStart=/root/start-ckan.sh
Restart=always

[Install]

WantedBy=default.target
